# sipXpbx/doc/developer/interop.pingtel.com/scripts/crontab
# Summarize traffic through the proxy, just after the logs are rotated.
30 4 * * *	$HOME/bin/proxy-summary
35 4 * * *	$HOME/bin/dialog-event-summary
